# .config/autostart/startEmacsAtLogin.sh.desktop
[Desktop Entry]
Exec=/home/luke/startEmacsAtLogin.sh
Icon=dialog-scripts
Name=startEmacsAtLogin.sh
Type=Application
X-KDE-AutostartScript=true
